# configs/archstrike/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --login-options "-f root" --noissue --noclear %I 38400 linux
